# main/supervisor.ini
[program:pyempty.web]
command=/root/.monitor/bin/gunicorn mainweb:app -b 0.0.0.0:8888 -w 4 -k gevent
directory=/root/github/pyempty/main
environment=PATH="/root/.monitor/bin",PYTHONPATH="../src"
stderr_logfile=/data/log/%(program_name)s.stderr.log
